# old_work/svn/gems/trunk/ecloud/artifacts/ecloud_3.6.x/master/14.04/DB1/var/spool/cron/crontabs/enlighted
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
*/5 * * * * php /home/enlighted/utils/generateSPPACerts.php >> /home/enlighted/utils/gencerts.log 2>&1 
